# systemd/redeem.service
[Unit]
Description=The Replicape Dameon

[Service]
ExecStart=/usr/src/redeem/software/Redeem.py

[Install]
WantedBy=basic.target
