# systemd/redeem.service
[Unit]
Description=The Replicape Dameon

[Service]
ExecStart=/usr/bin/redeem

[Install]
WantedBy=basic.target
