# debian/postinst
#!/bin/sh

# Enable and start redeem.service
systemctl enable redeem.service
systemctl start redeem.service
