# debian/prerm
#!/bin/sh

# Shutdown and remove redeem.service
systemctl stop redeem.service
systemctl disable redeem.service
